# root/etc/s6-overlay/s6-rc.d/svc-rustnzb/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid abc rustnzb \
    --config /config/config.toml \
    --data-dir /data \
    --port 9090
